# tools/nbd_backing/vramfs-nbd.service
[Unit]
Description=VRAM NBD export (prototype)
After=network.target

[Service]
Type=simple
User=root
ExecStart=/home/aman/Desktop/vramfs/bin/start_gpu_swap
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
[Unit]
Description=vramfs NBD server (GPU-backed)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nbd_server 10809
Restart=on-failure
RestartSec=5
WatchdogSec=30
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
